# excalibur/patchfolder/json2file-go/DEBIAN/postinst
#!/bin/sh
set -e
### see: dh_installdeb(1)
###      https://www.debian.org/doc/manuals/maint-guide/dother.en.html#maintscripts
### summary of how this script can be called:
###        * <postinst> `configure' <most-recently-configured-version>
###        * <old-postinst> `abort-upgrade' <new version>
###        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
###          <new-version>
###        * <postinst> `abort-remove'
###        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
###          <failed-install-package> <version> `removing'
###          <conflicting-package> <version>
### for details, see https://www.debian.org/doc/debian-policy/ or
### the debian-policy package

USER="www-data"
GROUP="www-data"
CONF_DIR="/var/spool/json2file-go"
case "$1" in
    configure)
		# Allow the daemon to bind to low level ports adjusting capabilities
		setcap CAP_NET_BIND_SERVICE=+eip /usr/sbin/json2file-go
		# Adjust directory permissions and ownerships using dpkg-statoverride
		if ! dpkg-statoverride --list ${CONF_DIR} >/dev/null 2>&1; then 
			dpkg-statoverride --update --add ${USER} ${GROUP} 750 ${CONF_DIR}
		fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

### dh_installdeb will replace this with shell code automatically
### generated by other debhelper scripts.

# Automatically added by dh_installtmpfiles/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -x "$(command -v systemd-tmpfiles)" ]; then
		systemd-tmpfiles ${DPKG_ROOT:+--root="$DPKG_ROOT"} --create json2file-go.conf || true
	fi
fi
# End automatically added section
# Automatically added by dh_installinit/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -z "$DPKG_ROOT" ] && [ -x "/etc/init.d/json2file-go" ]; then
		update-rc.d json2file-go defaults >/dev/null
		if [ -n "$2" ]; then
			_dh_action=restart
		else
			_dh_action=start
		fi
		invoke-rc.d --skip-systemd-native json2file-go $_dh_action || exit 1
	fi
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	# The following line should be removed in trixie or trixie+1
	deb-systemd-helper unmask 'json2file-go.service' >/dev/null || true

	# was-enabled defaults to true, so new installations run enable.
	if deb-systemd-helper --quiet was-enabled 'json2file-go.service'; then
		# Enables the unit on first installation, creates new
		# symlinks on upgrades if the unit file has changed.
		deb-systemd-helper enable 'json2file-go.service' >/dev/null || true
	else
		# Update the statefile to add new symlinks (if any), which need to be
		# cleaned up on purge. Also remove old symlinks.
		deb-systemd-helper update-state 'json2file-go.service' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -d /run/systemd/system ]; then
		systemctl --system daemon-reload >/dev/null || true
		if [ -n "$2" ]; then
			_dh_action=restart
		else
			_dh_action=start
		fi
		deb-systemd-invoke $_dh_action 'json2file-go.service' 'json2file-go.socket' >/dev/null || true
	fi
fi
# End automatically added section


exit 0
