# excalibur/patchfolder/lomiri-indicator-network/usr/lib/systemd/user/lomiri-indicator-network-secret-agent.service
[Unit]
Description=Indicator Network Secret Agent
PartOf=lomiri-indicator-network.service

[Service]
ExecStart=/usr/libexec/lomiri-indicator-network/lomiri-indicator-network-secret-agent
Restart=on-failure
