# excalibur/patchfolder/dnscache/etc/sv/dnscache/run
#!/usr/bin/env /lib/runit/invoke-run
exec 2>&1
exec </dev/urandom
exec chpst -Udnscache -o250 -d "${DATALIMIT}" /usr/sbin/dnscache
