# excalibur/patchfolder/walldns/etc/sv/walldns/run
#!/usr/bin/env /lib/runit/invoke-run
exec 2>&1
exec chpst -U walldns /usr/sbin/walldns
