# excalibur/patchfolder/apacheds/DEBIAN/prerm
#!/bin/sh
set -e
# Automatically added by dh_installsystemd/13.15.3
if [ -z "${DPKG_ROOT:-}" ] && [ "$1" = remove ] && [ -d /run/systemd/system ] ; then
	deb-systemd-invoke stop 'apacheds.service' >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installinit/13.15.3
if [ -z "${DPKG_ROOT:-}" ] && [ "$1" = remove ] && [ -x "/etc/init.d/apacheds" ] ; then
	invoke-rc.d --skip-systemd-native apacheds stop || exit 1
fi
# End automatically added section
