# excalibur/patchfolder/mobile-tweaks-common/usr/lib/systemd/system/ModemManager.service.d/mobile-tweaks.conf
# Use the experimental "quick suspend/resume" feature so the
# modem is available right after resuming from suspend; this
# avoids the (rather long) re-enumeration time, which often
# causes voice calls to be missed.
[Service]
ExecStart=
ExecStart=/usr/sbin/ModemManager --test-quick-suspend-resume
