# excalibur/patchfolder/vitrage-graph/DEBIAN/postrm
#!/bin/sh
set -e
# Automatically added by dh_installinit/13.24.2
if [ "$1" = "remove" ] && [ -x "/etc/init.d/vitrage-graph" ] ; then
	chmod -x "/etc/init.d/vitrage-graph" >/dev/null || true
fi
if [ -z "$DPKG_ROOT" ] && [ "$1" = "purge" ] ; then
	update-rc.d vitrage-graph remove >/dev/null
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "purge" ]; then
	if [ -x "/usr/bin/deb-systemd-helper" ]; then
		deb-systemd-helper purge 'vitrage-graph.service' >/dev/null || true
	fi
fi
# End automatically added section
