# excalibur/patchfolder/nextepc-pgw/usr/lib/systemd/system/nextepc-pgw.service
[Unit]
Description=Packet Data Network Gateway (EPC / LTE)

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/nextepc-pgwd -f /etc/nextepc/pgw.conf
RestartSec=2

[Install]
WantedBy=multi-user.target
