# excalibur/patchfolder/rpki-client/usr/lib/systemd/system/rpki-client.timer
[Unit]
Description=Periodically run rpki-client

[Timer]
OnBootSec=2min
OnUnitInactiveSec=15min

[Install]
WantedBy=timers.target

