# excalibur/patchfolder/runit-services/usr/share/runit/sv/openntpd/run
#!/usr/bin/env /lib/runit/invoke-run
#Copyright: 2022 Lorenzo Puliti <plorenzo@disroot.org>
#License: CC0-1.0

touch /etc/runit/override-sysv.d/${PWD##*/}.pkgrunit

if [ ! -d "/run/openntpd" ]; then 
    mkdir -p /run/openntpd
fi
/usr/sbin/openntpd -n  2>&1 || exit 162

exec 2>&1

exec ${XCHPST:-} /usr/sbin/openntpd -d -s $DAEMON_OPTS
