# excalibur/patchfolder/runit-services/usr/share/runit/sv/exim4/run
#!/usr/bin/env /lib/runit/invoke-run
#Copyright: 2022 Lorenzo Puliti <plorenzo@disroot.org>
#           2024-2025 Andrew Bower <andrew@bower.uk>
#License: CC0-1.0

exec 2>&1

# Compatibility with defaults file carried over from pre-trixie installations
# this compat will be removed in future - likely in trixie +1
if [ -n "$QUEUERUNNER$QUEUEINTERVAL$COMMONOPTIONS$QUEUERUNNEROPTIONS$QFLAGS$SMTPLISTENEROPTIONS" ]
then
  EXIMSERVICE=${EXIMSERVICE--bdf -q${QFLAGS}${QUEUEINTERVAL} $COMMONOPTIONS $QUEUERUNNEROPTIONS $SMTPLISTENEROPTIONS}
else
  EXIMSERVICE=${EXIMSERVICE--bdf -q30m}
fi

[ -d /run/exim4 ] || \
  install -d -oDebian-exim -gDebian-exim -m750 /run/exim4
LANG=C
export LANG

# '-d' option appended to provide debug output for svlogd

/usr/sbin/update-exim4.conf $UPEX4OPTS || exit 1
exec ${XCHPST:-} /usr/sbin/exim4 $EXIMSERVICE -d
