# excalibur/patchfolder/runit-services/usr/share/runit/sv/network-manager/run
#!/usr/bin/env /lib/runit/invoke-run
#Copyright: 2023 Lorenzo Puliti <plorenzo@disroot.org>
#License: CC0-1.0

exec 2>&1

sv start dbus || exit 170

exec ${XCHPST:-chpst} -o 65536 /usr/sbin/NetworkManager --no-daemon
