# excalibur/patchfolder/oar-server/usr/lib/systemd/system/oar-accounting.timer
[Unit]
Description=OAR accounting timer
Documentation=man:oaraccounting(8)

[Timer]
OnCalendar=*-*-* *:00:00

[Install]
WantedBy=timers.target
