# excalibur/patchfolder/oar-server/usr/lib/systemd/system/oar-accounting.service
[Unit]
Description=OAR accounting service
Documentation=man:oar-server(8)
After=oar-server.service

[Service]
User=root
ExecStart=/usr/sbin/oaraccounting
