# excalibur/patchfolder/painintheapt/usr/lib/systemd/system/painintheapt-weekly.timer
[Unit]
Description=Weekly painintheapt

[Timer]
OnCalendar=Sun
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
