# excalibur/patchfolder/runit-services/usr/share/runit/sv/chrony/run
#!/usr/bin/env /lib/runit/invoke-run
#Copyright: 2005-2008 Gerrit Pape <pape@smarden.org>
#           2022 Lorenzo Puliti <plorenzo@disroot.org>
#License: BSD-3-Clause

exec 2>&1

exec /usr/sbin/chronyd -d $DAEMON_OPTS
