# excalibur/patchfolder/mythtv-status/usr/lib/systemd/system/mythtv-status_update-motd.service
[Unit]
Description=mythtv-status MOTD update
Documentation=man:mythtv-update-motd(8)

[Service]
Type=oneshot
ExecStart=/usr/sbin/mythtv-update-motd
