# excalibur/patchfolder/qrtr-tools/usr/lib/systemd/system/qrtr-ns.service
[Unit]
Description=QIPCRTR Name Service

[Service]
ExecStart=/usr/bin/qrtr-ns -f 1
Restart=always

[Install]
WantedBy=multi-user.target
