# excalibur/patchfolder/pagure/DEBIAN/postinst
#!/bin/sh
# postinst script for pagure
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

# Source debconf library.
if test -f /usr/share/debconf/confmodule ; then
    . /usr/share/debconf/confmodule
fi

# Source the information about the pagure user.
PAGURE_USER="_pagure"
PAGURE_HOME="/var/lib/pagure"

# Helper function to create a directory and perform chown on it.
my_mkdir ()
{
    d=$1
    mkdir -p $d
    chown ${PAGURE_USER}:${PAGURE_USER} $d
}

case "$1" in
    configure)
	# Create the pagure user.
	if ! getent passwd "$PAGURE_USER" > /dev/null; then
	    adduser --quiet --gecos 'pagure user' \
		    --home "$PAGURE_HOME" --disabled-password \
		    --force-badname "$PAGURE_USER"
	fi

	# Create the directory tree.
	my_mkdir /var/www/releases

	my_mkdir ${PAGURE_HOME}/repositories
	my_mkdir ${PAGURE_HOME}/repositories/docs
	my_mkdir ${PAGURE_HOME}/repositories/forks
	my_mkdir ${PAGURE_HOME}/repositories/requests
	my_mkdir ${PAGURE_HOME}/repositories/tickets

	my_mkdir ${PAGURE_HOME}/remotes

	my_mkdir ${PAGURE_HOME}/.ssh
	chmod 700 ${PAGURE_HOME}/.ssh

	my_mkdir ${PAGURE_HOME}/.gitolite
	my_mkdir ${PAGURE_HOME}/.gitolite/conf
	my_mkdir ${PAGURE_HOME}/.gitolite/keydir
	my_mkdir ${PAGURE_HOME}/.gitolite/logs
	touch ${PAGURE_HOME}/.gitolite/conf/gitolite.conf
	cp -a /usr/share/pagure/gitolite3.rc ${PAGURE_HOME}/.gitolite.rc
	chown ${PAGURE_USER}:${PAGURE_USER} \
	      ${PAGURE_HOME}/.gitolite/conf/gitolite.conf \
	      ${PAGURE_HOME}/.gitolite.rc

	echo
	echo "*** Pagure needs to be configured by hand                      ***"
	echo "*** Please, take a look at /usr/share/doc/pagure/README.Debian ***"
	echo "*** and follow the instructions there.                         ***"
	echo
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.


# Automatically added by dh_python3
if command -v py3compile >/dev/null 2>&1; then
	py3compile -p pagure 
fi
if command -v pypy3compile >/dev/null 2>&1; then
	pypy3compile -p pagure  || true
fi

if command -v py3compile >/dev/null 2>&1; then
	py3compile -p pagure /usr/share/pagure
fi
if command -v pypy3compile >/dev/null 2>&1; then
	pypy3compile -p pagure /usr/share/pagure || true
fi

# End automatically added section
# Automatically added by dh_installsystemd/13.24.1
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'pagure-worker.service'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'pagure-worker.service' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'pagure-worker.service'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'pagure-worker.service' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'pagure-worker.service' >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.1
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'pagure-gitolite-worker.service'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'pagure-gitolite-worker.service' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'pagure-gitolite-worker.service'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'pagure-gitolite-worker.service' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'pagure-gitolite-worker.service' >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.1
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'pagure-api-key-expire-mail.timer'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'pagure-api-key-expire-mail.timer' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'pagure-api-key-expire-mail.timer'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'pagure-api-key-expire-mail.timer' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'pagure-api-key-expire-mail.timer' >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.1
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'pagure-mirror-project-in.timer'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'pagure-mirror-project-in.timer' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'pagure-mirror-project-in.timer'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'pagure-mirror-project-in.timer' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'pagure-mirror-project-in.timer' >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.1
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'pagure-authorized-keys-worker.service'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'pagure-authorized-keys-worker.service' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'pagure-authorized-keys-worker.service'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'pagure-authorized-keys-worker.service' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'pagure-authorized-keys-worker.service' >/dev/null || true
fi
# End automatically added section


exit 0
