# excalibur/patchfolder/runit-services/usr/share/runit/sv/anacron/run
#!/usr/bin/env /lib/runit/invoke-run
#Copyright: 2022 Lorenzo Puliti <plorenzo@disroot.org>
#License: CC0-1.0

touch /etc/runit/override-sysv.d/${PWD##*/}.pkgrunit

# exit status of on_ac_power:
# 0= on-ac // 1= on-battery // 255=unknown, likely a desktop without APM
if [ x"$ANACRON_RUN_ON_BATTERY_POWER" != x"yes" ] && [ -x /usr/bin/on_ac_power ]; then
    /usr/bin/on_ac_power || retval=$?
    if [ x"$retval" = x1 ]; then
        echo "deferred while on battery power" && exit 161
    fi
fi
exec 2>&1

# don't restart anacron when it's done
sv once anacron

exec /usr/sbin/anacron -d $ANACRON_ARGS
