# excalibur/patchfolder/partimage-server/DEBIAN/postrm
#!/bin/sh

set -e

if [ "$1" = "purge" ]; then
	rm -f /etc/partimaged/partimaged.cert \
		/etc/partimaged/partimaged.key \
		/etc/partimaged/privkey.pem \
		/etc/partimaged/partimagedusers \
		/etc/partimaged/partimaged.rand \
		/etc/partimaged/passwd.db

	# I'm not using "rm -rf" due the user could have additional data
	rmdir /etc/partimaged || true
fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

# Automatically added by dh_installdeb/13.14.1
dpkg-maintscript-helper rm_conffile /etc/init.d/partimaged 0.6.9-10\~ -- "$@"
dpkg-maintscript-helper rm_conffile /etc/default/partimaged 0.6.9-10\~ -- "$@"
# End automatically added section
# Automatically added by dh_installsystemd/13.14.1
if [ "$1" = remove ] && [ -d /run/systemd/system ] ; then
	systemctl --system daemon-reload >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.14.1
if [ "$1" = "purge" ]; then
	if [ -x "/usr/bin/deb-systemd-helper" ]; then
		deb-systemd-helper purge 'partimaged.service' >/dev/null || true
	fi
fi
# End automatically added section


