# excalibur/patchfolder/kio-fuse/usr/lib/systemd/user/kio-fuse.service
[Unit]
Description=Fuse interface for KIO
PartOf=graphical-session.target

[Service]
ExecStart=/usr/lib/x86_64-linux-gnu/libexec/kio-fuse -f
BusName=org.kde.KIOFuse
Slice=background.slice
