# excalibur/patchfolder/restic-rest-server/DEBIAN/postinst
#!/bin/sh

set -e

# summary of how this script can be called:
#        * <postinst> 'configure' <most-recently-configured-version>
#        * <old-postinst> 'abort-upgrade' <new version>
#        * <conflictor's-postinst> 'abort-remove' 'in-favour' <package>
#          <new-version>
#        * <deconfigured's-postinst> 'abort-deconfigure' 'in-favour'
#          <failed-install-package> <version> 'removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package
#
# quoting from the policy:
#     Any necessary prompting should almost always be confined to the
#     post-installation script, and should be protected with a conditional
#     so that unnecessary prompting doesn't happen if a package's
#     installation fails and the 'postinst' is called with 'abort-upgrade',
#     'abort-remove' or 'abort-deconfigure'.

case "$1" in
  configure)
    # Add restic-rest-server user
    adduser --quiet --system --home /nonexistent --no-create-home \
            --group --gecos "Restic REST server" restic-rest-server

    CFG=/etc/restic-rest-server/users.htpasswd
    if ! dpkg-statoverride --list "$CFG" >/dev/null;
    then
        chown root:restic-rest-server "$CFG"
        chmod 640 "$CFG"
    fi
  ;;

  abort-upgrade|abort-remove|abort-deconfigure)
    :
  ;;

  *)
    echo "postinst called with unknown argument '$1'" >&2
    exit 1
  ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'restic-rest-server.service'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'restic-rest-server.service' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'restic-rest-server.service'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'restic-rest-server.service' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'restic-rest-server.service' >/dev/null || true
fi
# End automatically added section


exit 0
