# excalibur/patchfolder/aodh-api/usr/lib/systemd/system/aodh-api.service
[Unit]
Description=OpenStack Aodh API (aodh-api)
After=mongodb.service rabbitmq-server.service frr.service network-online.target local-fs.target remote-fs.target 
Wants=mongodb.service rabbitmq-server.service frr.service network-online.target 

Documentation=man:aodh-api(1)

[Service]
User=aodh
Group=aodh
Type=simple
WorkingDirectory=~
RuntimeDirectory=aodh lock/aodh
CacheDirectory=aodh
ExecStart=/etc/init.d/aodh-api systemd-start
Restart=on-failure
LimitNOFILE=65535
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
