# excalibur/patchfolder/lomiri-history-service/usr/lib/systemd/user/lomiri-history-daemon.service
[Unit]
Description=History service to store messages and calls

[Service]
Type=dbus
BusName=com.lomiri.HistoryService
ExecStart=/usr/libexec/lomiri-history-daemon
