# excalibur/patchfolder/librem5-tweaks/usr/lib/systemd/system/gpsd.service.d/librem5.conf
[Service]
EnvironmentFile=
EnvironmentFile=/etc/default/librem5-gpsd

