# excalibur/patchfolder/x2gothinclient-displaymanager/DEBIAN/prerm
#!/bin/sh

# prerm script for x2gothinclient-displaymanager
#
# see: dh_installdeb(1)

set -e

# source debconf's configuration module
. /usr/share/debconf/confmodule

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#        <package-being-installed> <version> `removing'
#        <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

THIS_PACKAGE=x2gothinclient-displaymanager
DEFAULT_DISPLAY_MANAGER_FILE=/etc/X11/default-display-manager

if [ "$1" = "remove" -o "$1" = "deconfigure" ]; then

	# disown this question
	db_unregister shared/default-x-display-manager || true

	# does the question still exist?
	if db_get shared/default-x-display-manager; then
		db_metaget shared/default-x-display-manager owners
		db_subst shared/default-x-display-manager choices "$RET"
		db_get shared/default-x-display-manager
		# are we removing the currently selected display manager?
		if [ "$THIS_PACKAGE" = "$RET" ]; then
			if [ -e "$DEFAULT_DISPLAY_MANAGER_FILE" ]; then
				if [ "$(cat $DEFAULT_DISPLAY_MANAGER_FILE)" = "/usr/sbin/x2gothinclientd" ]; then
					rm "$DEFAULT_DISPLAY_MANAGER_FILE"
				fi
			fi

			# ask the user to choose a new default
			db_fset shared/default-x-display-manager seen false
			db_input critical shared/default-x-display-manager || true
			db_go

			# if the display manager file doesn't exist, write it with the path
			# to the new default display manager
			if [ ! -e $DEFAULT_DISPLAY_MANAGER_FILE ]; then
				db_get shared/default-x-display-manager
				echo "Please be sure to run \"dpkg-reconfigure $RET\"."
				db_get "$RET"/daemon_name
				echo "$RET" > "$DEFAULT_DISPLAY_MANAGER_FILE"
			fi
		fi
	fi
fi


case "$1" in
	remove)
		:
		;;
	deconfigure|upgrade|failed-upgrade)
		:
		;;
	*) echo "$0: didn't understand being called with \`$1'" 1>&2
		exit 1
		;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

# Automatically added by dh_installsystemd/13.24.1
if [ -z "$DPKG_ROOT" ] && [ "$1" = remove ] && [ -d /run/systemd/system ] ; then
	deb-systemd-invoke stop 'x2gothinclient-displaymanager.service' >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installinit/13.24.1
if [ -z "$DPKG_ROOT" ] && [ "$1" = remove ] && [ -x "/etc/init.d/x2gothinclient-displaymanager" ] ; then
	invoke-rc.d --skip-systemd-native x2gothinclient-displaymanager stop || exit 1
fi
# End automatically added section


exit 0
