# excalibur/patchfolder/hydrapaper/usr/lib/systemd/user/org.gabmus.hydrapaperd.service
[Unit]
Description=HydraPaper Daemon

[Service]
Type=dbus
BusName=org.gabmus.hydrapaper.Daemon
ExecStart=/usr/libexec/hydrapaperd
