# excalibur/patchfolder/runit-services/usr/share/runit/sv/mini-httpd/run
#!/usr/bin/env /lib/runit/invoke-run
#Copyright: 2017 Dmitry Bogatov <KAction@disroot.org>
#             2023  Lorenzo Puliti <plorenzo@disroot.org>
#License: GPL-3+

exec 2>&1

exec /usr/sbin/mini_httpd -C /etc/mini-httpd.conf  -D -l /dev/stdout
