# excalibur/patchfolder/ltspfsd/lib/systemd/system/ltspfsd.service
[Unit]
Description=LTSP Filesystem Daemon

[Service]
Type=forking
ExecStart=/usr/bin/ltspfsd
