# excalibur/patchfolder/lomiri-download-manager/usr/lib/systemd/user/lomiri-download-manager-systemd.service
[Unit]
Description=Lomiri Download Manager Service

[Service]
Type=dbus
BusName=com.lomiri.applications.Downloader
ExecStart=/usr/bin/lomiri-download-manager
