# excalibur/patchfolder/socklog-run/etc/sv/socklog-klog/run
#!/bin/sh
exec 2>&1
exec </proc/kmsg
exec chpst -u _socklog-klog socklog ucspi
