# excalibur/patchfolder/bcron/etc/sv/bcron-spool/run
#!/bin/sh
exec 2>&1
sv up bcron-update
exec chpst -Ucron unixserver -U "${BCRON_SOCKET:-/var/run/bcron-spool}" bcron-spool
