# excalibur/patchfolder/runit-services/usr/share/runit/sv/slim/run
#!/usr/bin/env /lib/runit/invoke-run
#Copyright: 2022 Lorenzo Puliti <plorenzo@disroot.org>
#License: CC0-1.0

DEFAULT_DISPLAY_MANAGER=$(cat /etc/X11/default-display-manager)

# don't start if slim is not the default DM
if [ $DEFAULT_DISPLAY_MANAGER != /usr/bin/slim ]; then
	exit 161
fi

sv start dbus || true
sv start elogind || true

if [ -x /bin/plymouth ] ; then
        /bin/plymouth quit || true
fi

exec 2>&1

exec /usr/bin/slim -n
