# excalibur/patchfolder/sms4you/usr/lib/systemd/system/sms4you.service
[Unit]
Description=sms4you daemon
Wants=ModemManager.service prosody.service

[Service]
ExecStart=/usr/bin/sms4you
Restart=always
Type=notify
User=_sms4you

[Install]
WantedBy=multi-user.target
