# excalibur/patchfolder/runit-services/usr/share/runit/sv/atd/run
#!/usr/bin/env /lib/runit/invoke-run
#Copyright: 2022 Lorenzo Puliti <plorenzo@disroot.org>
#License: CC0-1.0

sv start default-syslog ||  true

exec 2>&1

exec /usr/sbin/atd -f
