# excalibur/patchfolder/runit-services/usr/share/runit/sv/unbound/run
#!/usr/bin/env /lib/runit/invoke-run
#Copyright: 2024 Dimitris T. <dimitris@stinpriza.org>
#License: BSD-3-Clause

exec 2>&1

if ! unbound-checkconf ; then
        echo "wrong syntax in confg file"
        sv d .
fi

# TODO: running in chroot not supported for now

exec /usr/sbin/unbound -d -d -p
