# excalibur/patchfolder/rancid/DEBIAN/postinst
#! /bin/sh
# postinst script for rancid
#
# see: dh_installdeb(1)

set -e

. /usr/share/debconf/confmodule
db_version 2.0

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package
#
# quoting from the policy:
#     Any necessary prompting should almost always be confined to the
#     post-installation script, and should be protected with a conditional
#     so that unnecessary prompting doesn't happen if a package's
#     installation fails and the `postinst' is called with `abort-upgrade',
#     `abort-remove' or `abort-deconfigure'.



case "$1" in
    configure)
        # change wrong shell /bin/false for user rancid to /bin/bash:
	case "`getent passwd rancid`" in
	    *:/bin/false)
	    usermod -s /bin/bash rancid
	    ;;
	esac

	[ "$RET" = "false" ] && exit 0
	if ! `getent passwd rancid >/dev/null`; then
		adduser --system --group --home /var/lib/rancid --quiet \
		    --shell /bin/bash rancid
        fi
	if [ -d /var/log/rancid ]; then
		chown rancid:rancid /var/log/rancid
	fi
	if [ -d /var/lib/rancid ]; then
		chown rancid:rancid /var/lib/rancid
	fi
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)

    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac


db_stop

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'rancid.timer'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'rancid.timer' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'rancid.timer'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'rancid.timer' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'rancid.timer' >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -d /run/systemd/system ]; then
		systemctl --system daemon-reload >/dev/null || true
		if [ -n "$2" ]; then
			_dh_action=restart
		else
			_dh_action=start
		fi
		deb-systemd-invoke $_dh_action 'rancid.service' 'rancid.timer' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'rancid-cleanup.timer'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'rancid-cleanup.timer' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'rancid-cleanup.timer'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'rancid-cleanup.timer' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'rancid-cleanup.timer' >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -d /run/systemd/system ]; then
		systemctl --system daemon-reload >/dev/null || true
		if [ -n "$2" ]; then
			_dh_action=restart
		else
			_dh_action=start
		fi
		deb-systemd-invoke $_dh_action 'rancid-cleanup.service' 'rancid-cleanup.timer' >/dev/null || true
	fi
fi
# End automatically added section


exit 0
