# excalibur/patchfolder/sitesummary/DEBIAN/postrm
#!/bin/sh

set -e

if [ "$1" = "purge" ] ; then
    cd /var/lib/sitesummary
    rm -rf entries tmpstorage www munin-conf.d
    rm -f munin.conf
    rm -f nagios-generated.cfg
fi

# Automatically added by dh_installsystemd/13.15.3
if [ "$1" = remove ] && [ -d /run/systemd/system ] ; then
	systemctl --system daemon-reload >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.15.3
if [ "$1" = "purge" ]; then
	if [ -x "/usr/bin/deb-systemd-helper" ]; then
		deb-systemd-helper purge 'sitesummary.timer' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_installdebconf/13.15.3
if [ "$1" = purge ] && [ -e /usr/share/debconf/confmodule ]; then
	. /usr/share/debconf/confmodule
	db_purge
fi
# End automatically added section

