# excalibur/patchfolder/chkboot/DEBIAN/postinst
#!/bin/sh
# postinst script for chkboot

set -e

case "$1" in
    configure)
        # On upgrade
        if [ -n "$2" ]; then

            # For version <= 1.3-2
            if dpkg --compare-versions "$2" le "1.3-2"; then

                # Update permission to match package one.
                # See #192981
                if [ -x /etc/profile.d/chkboot-profilealert.sh ]; then
                    chmod 644 /etc/profile.d/chkboot-profilealert.sh
                fi

            fi

        fi
    ;;

    triggered)
        # This file trigger an update in the apt hook
        touch "/var/lib/chkboot/needs-update" || true
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

# Automatically added by dh_installdeb/13.15.3
dpkg-maintscript-helper rm_conffile /etc/kernel/postinst.d/zzz-chkboot 1.3-4\~ -- "$@"
dpkg-maintscript-helper rm_conffile /etc/kernel/postrm.d/zzz-chkboot 1.3-4\~ -- "$@"
# End automatically added section
# Automatically added by dh_installinit/13.15.3
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -z "${DPKG_ROOT:-}" ] && [ -x "/etc/init.d/chkboot" ]; then
		update-rc.d chkboot defaults >/dev/null
		if [ -n "$2" ]; then
			_dh_action=restart
		else
			_dh_action=start
		fi
		invoke-rc.d --skip-systemd-native chkboot $_dh_action || exit 1
	fi
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.15.3
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	# The following line should be removed in trixie or trixie+1
	deb-systemd-helper unmask 'chkboot.service' >/dev/null || true

	# was-enabled defaults to true, so new installations run enable.
	if deb-systemd-helper --quiet was-enabled 'chkboot.service'; then
		# Enables the unit on first installation, creates new
		# symlinks on upgrades if the unit file has changed.
		deb-systemd-helper enable 'chkboot.service' >/dev/null || true
	else
		# Update the statefile to add new symlinks (if any), which need to be
		# cleaned up on purge. Also remove old symlinks.
		deb-systemd-helper update-state 'chkboot.service' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.15.3
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -d /run/systemd/system ]; then
		systemctl --system daemon-reload >/dev/null || true
		if [ -n "$2" ]; then
			_dh_action=restart
		else
			_dh_action=start
		fi
		deb-systemd-invoke $_dh_action 'chkboot.service' >/dev/null || true
	fi
fi
# End automatically added section


# On upgrade, after rm_conffiles, remove emtpy directories
# See #665797
if [ "$1" = "configure" ]; then
    if [ -n "$2" ]; then
        if dpkg --compare-versions "$2" le "1.3-4"; then
            if [ -d /etc/kernel/postinst.d ]; then
                rmdir --ignore-fail-on-non-empty /etc/kernel/postinst.d
            fi
            if [ -d /etc/kernel/postrm.d ]; then
                rmdir --ignore-fail-on-non-empty /etc/kernel/postrm.d
            fi
        fi
    fi
fi

exit 0
