# excalibur/patchfolder/dmarcts-report-parser/usr/lib/systemd/system/dmarcts-report-parser.timer
[Unit]
Description=Run dmarcts-report-parser once per day.

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
