# excalibur/patchfolder/openhpid/DEBIAN/postrm
#!/bin/sh
set -e
# Automatically added by dh_installinit/13.24.1
if [ "$1" = "remove" ] && [ -x "/etc/init.d/openhpid" ] ; then
	chmod -x "/etc/init.d/openhpid" >/dev/null || true
fi
if [ -z "$DPKG_ROOT" ] && [ "$1" = "purge" ] ; then
	update-rc.d openhpid remove >/dev/null
fi
# End automatically added section
