# excalibur/patchfolder/runit-services/usr/share/runit/sv/lldpd/run
#!/usr/bin/env /lib/runit/invoke-run
#Copyright: 2024 Andrew Bower <andrew@bower.uk>
#License: BSD-3-Clause

exec 2>&1

# Start daemon. See lldpd(8):
#  - 1st '-d' prevents autodaemonisation
#  - 2nd '-d' logs to stderr rather than syslog
#  - 3rd '-d' enables verbose logging
#  - a 4th '-d' would enable debug logging
exec /usr/sbin/lldpd -d -d -d $DAEMON_ARGS
