# excalibur/patchfolder/spamassassin/usr/lib/systemd/system/spamassassin-maintenance.timer
[Unit]
Description=Daily spamassassin maintenance

[Timer]
OnCalendar=*-*-* 6:00
RandomizedDelaySec=60m
Persistent=true

[Install]
WantedBy=timers.target
