# excalibur/patchfolder/runit-services/usr/share/runit/sv/lightdm/run
#!/usr/bin/env /lib/runit/invoke-run
#Copyright: 2022 Lorenzo Puliti <plorenzo@disroot.org>
#License: CC0-1.0

DEFAULT_DISPLAY_MANAGER=$(cat /etc/X11/default-display-manager)
# start only if lightdm the default DM
if [ $DEFAULT_DISPLAY_MANAGER != /usr/sbin/lightdm ]; then
	exit 161
fi

sv start dbus || exit 170
sv start elogind || true

if [ -r /etc/default/locale ]; then
  . /etc/default/locale
  export LANG LANGUAGE
fi

exec 2>&1

exec /usr/sbin/lightdm 
