# excalibur/patchfolder/runit-services/usr/share/runit/sv/isc-dhcp-server.6/run
#!/usr/bin/env /lib/runit/invoke-run
#Copyright: 2005-2008 Gerrit Pape <pape@smarden.org>
#        2023 Lorenzo Puliti <plorenzo@disroot.org>
#        2024 Andrew Bower <andrew@bower.uk>
#License: BSD-3-Clause

exec 2>&1

if [ -z "$INTERFACESv6" ]; then
        [ -r /etc/default/isc-dhcp-server ] && . /etc/default/isc-dhcp-server
        [ -n "$INTERFACESv6" ] || {
          echo "invoke-run: no IPv6 interfaces specified for ${PWD##*/}"
          exit 162
        }
fi

leases="/var/lib/dhcp/dhcpd6.leases"
[ -e "$leases" ] || touch "$leases"

/usr/sbin/dhcpd -6 -f -d -t $INTERFACESv6 || exit 162
/usr/sbin/dhcpd -6 -f -d -T $INTERFACESv6 || exit 162

exec /usr/sbin/dhcpd -6 -f -d $INTERFACESv6
