# excalibur/patchfolder/buteo-syncfw/usr/lib/systemd/user/msyncd.service
[Unit]
Description=Sync FW daemon

[Service]
Type=dbus
ExecStart=/usr/bin/msyncd
BusName=com.meego.msyncd
