# excalibur/patchfolder/arno-iptables-firewall/DEBIAN/postrm
#!/bin/sh
# vim: syntax=sh ts=4 sw=4 sts=4 sr noet
# postrm script for arno-iptables-firewall
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postrm> `remove'
#        * <postrm> `purge'
#        * <old-postrm> `upgrade' <new-version>
#        * <new-postrm> `failed-upgrade' <old-version>
#        * <new-postrm> `abort-install'
#        * <new-postrm> `abort-install' <old-version>
#        * <new-postrm> `abort-upgrade' <old-version>
#        * <disappearer's-postrm> `disappear' <overwriter>
#          <overwriter-version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package


# start up debconf here. 
# why? see https://manpages.debian.org/testing/debconf-doc/debconf-devel.7.en.html, paragraph HACKS. 
# plus automatically inserted code may be using debconf
# (in fact it does, and loads confmodule a little bit late IMHO)
. /usr/share/debconf/confmodule
db_version 2.0

case "$1" in
	purge)
		# when purging remove debconf managed config file
		rm -f /etc/arno-iptables-firewall/conf.d/00debconf.conf
	;;
	
	remove|disappear)
		# Clean-up code regarding systemd has been injected below by
		# dh_installsystemd into the DEBHELPER section.
		# Regarding sysvinit we need to do it ourselves, dh_installinit
		# is apparently inactive on this.
		update-rc.d -f arno-iptables-firewall remove >/dev/null || true
	;;

	upgrade|failed-upgrade|abort-install|abort-upgrade)
		# we don't stop on upgrades so that the firewall keeps operational.
		# this works since there's only a set of iptables rules, no real daemon.
		# eventually postinst configure replaces them with a new set of rules.
	;;

	*)
		echo "postrm called with unknown argument '$1'" >&2
		exit 1
	;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

# Automatically added by dh_installsystemd/13.24.1
if [ "$1" = remove ] && [ -d /run/systemd/system ] ; then
	systemctl --system daemon-reload >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.1
if [ "$1" = "purge" ]; then
	if [ -x "/usr/bin/deb-systemd-helper" ]; then
		deb-systemd-helper purge 'arno-iptables-firewall.service' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_installdebconf/13.24.1
if [ "$1" = purge ] && [ -e /usr/share/debconf/confmodule ]; then
	. /usr/share/debconf/confmodule
	db_purge
fi
# End automatically added section


exit 0
