# excalibur/patchfolder/nordugrid-arc-arex/DEBIAN/prerm
#!/bin/sh

set -e

# Automatically added by dh_installsystemd/13.24.2
if [ -z "$DPKG_ROOT" ] && [ "$1" = remove ] && [ -d /run/systemd/system ] ; then
	deb-systemd-invoke stop 'arc-arex-ws.service' >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.2
if [ -z "$DPKG_ROOT" ] && [ "$1" = remove ] && [ -d /run/systemd/system ] ; then
	deb-systemd-invoke stop 'arc-arex.service' >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installinit/13.24.2
if [ -z "$DPKG_ROOT" ] && [ "$1" = remove ] && [ -x "/etc/init.d/arc-arex-ws" ] ; then
	invoke-rc.d --skip-systemd-native arc-arex-ws stop || exit 1
fi
# End automatically added section
# Automatically added by dh_installinit/13.24.2
if [ -z "$DPKG_ROOT" ] && [ "$1" = remove ] && [ -x "/etc/init.d/arc-arex" ] ; then
	invoke-rc.d --skip-systemd-native arc-arex stop || exit 1
fi
# End automatically added section

# Automatically added by dh_python3
if command -v py3clean >/dev/null 2>&1; then
	py3clean -p nordugrid-arc-arex:amd64 
else
	dpkg -L nordugrid-arc-arex:amd64 | sed -En -e '/^(.*)\/(.+)\.py$/s,,rm "\1/__pycache__/\2".*,e'
	find /usr/lib/python3/dist-packages/ -type d -name __pycache__ -empty -print0 | xargs --null --no-run-if-empty rmdir
fi

# End automatically added section


if [ "$1" = "remove" ] ; then
   # Don't regenerate the bytecode that was removed by the autogenerated
   # debhelper maintainer script above...
   PYTHONDONTWRITEBYTECODE=1 arcctl test-ca cleanup
   rmdir /etc/grid-security/certificates 2>/dev/null || :
   rmdir /etc/grid-security 2>/dev/null || :
fi
