# excalibur/patchfolder/json2file-go/etc/init.d/json2file-go
#!/usr/bin/env /lib/init/init-d-script
## Debug shebang
##!/usr/bin/env -S /bin/sh -x /lib/init/init-d-script
### BEGIN INIT INFO
# Provides:          json2file-go
# Required-Start:    $syslog $time $remote_fs
# Required-Stop:     $syslog $time $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: launch json2file-go
# Description:       Debian init script to start json2file-go
### END INIT INFO
# For this init file to work the sysvinit-utils package must be installed
DAEMON="/usr/sbin/json2file-go"
VARIABLES="/usr/lib/json2file-go/variables.sh"
PIDFILE="/run/json2file-go/pid"
LOGFILE="/var/log/json2file-go.log"
START_ARGS="--background --no-close --make-pidfile"
do_start_cmd_override() {
    if [ -f "$VARIABLES" ]; then
        . "$VARIABLES"
	fi
	echo ""
	exec 1>> "$LOGFILE" 2>&1
	do_start_cmd
}
