# excalibur/patchfolder/protection-domain-mapper/usr/lib/systemd/system/pd-mapper.service
[Unit]
Description=Qualcomm PD mapper service

[Service]
ExecStart=/usr/bin/pd-mapper
Restart=always

[Install]
WantedBy=multi-user.target
