# excalibur/patchfolder/runit-services/usr/share/runit/sv/postfix/finish
#!/bin/sh
set -e

#Restore default values
#main.cf
postconf maillog_file= 

exec /lib/runit/finish-exec "$@"
