# excalibur/patchfolder/jetty12/usr/lib/systemd/system/jetty12.timer
[Unit]
Description=Weekly trigger for jetty12 logging. 

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
