# excalibur/patchfolder/runit-services/usr/share/runit/sv/preload/run
#!/usr/bin/env /lib/runit/invoke-run
#Copyright: 2022 Lorenzo Puliti <plorenzo@disroot.org>
#License: CC0-1.0

exec 2>&1

exec /usr/sbin/preload -s /var/lib/preload/preload.state -f  -l ''
