# excalibur/patchfolder/mpd-sima/DEBIAN/postinst
#!/bin/sh
# postinst script for mpd-sima
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

do_mpd_sima_adduser () {
  # add group
  if ! getent group mpd-sima >/dev/null ; then
    addgroup --quiet --system mpd-sima
  fi
  # add user
  if ! getent passwd mpd-sima >/dev/null; then
    echo -n "Creating system user 'mpd-sima'"
    adduser --quiet --ingroup mpd-sima --system --no-create-home --home /var/lib/mpd-sima mpd-sima && echo "."
  fi
}

do_mpd_sima_perms () {
  if ! dpkg-statoverride --list --quiet "/var/log/mpd-sima" >/dev/null; then
    dpkg-statoverride --force-all --quiet --update --add mpd-sima mpd-sima 0755 "/var/log/mpd-sima"
  fi
  if ! dpkg-statoverride --list --quiet "/etc/mpd-sima.cfg" >/dev/null; then
    dpkg-statoverride --force-all --quiet --update --add mpd-sima mpd-sima 0600 "/etc/mpd-sima.cfg"
  fi
}

do_mpd_sima_create_db () {
  # mpd-sima group should have write access to database to ease users access
  local simalib=/var/lib/mpd-sima
  local simadb=/var/lib/mpd-sima/sima.db
  if ! dpkg-statoverride --list --quiet "${simalib}" >/dev/null; then
    dpkg-statoverride --force-all --quiet --update --add mpd-sima mpd-sima 0775 "${simalib}"
  fi
  # TODO: Remove when bookworm is released (upgrading package over mpd-sima >= 0.18.0)
  # Ensure the database is migrated when previous version is there already
  if [ -e "${simadb}" ] ; then
    dpkg-statoverride --force-statoverride-remove --quiet --remove "${simadb}"
    echo -n "Migrate database if needed"
    su -s /bin/sh -c "mpd-sima --config=/etc/mpd-sima.cfg --var-dir=${simalib} --log-level=warning create-db" mpd-sima > \
      /dev/null && echo "."
    dpkg-statoverride --force-statoverride-add --quiet --update --add mpd-sima mpd-sima 0664 "${simadb}"
  fi
  #/TODO
  if [ ! -e "${simadb}" ] ; then
    echo -n "Creating database"
    su -s /bin/sh -c "mpd-sima --config=/etc/mpd-sima.cfg --var-dir=${simalib} --log-level=warning create-db" mpd-sima > \
      /dev/null && echo "."
  fi
  if ! dpkg-statoverride --list --quiet "${simadb}" >/dev/null; then
    dpkg-statoverride --force-all --quiet --update --add mpd-sima mpd-sima 0664 "${simadb}"
  fi
}

case "$1" in
  configure)
    do_mpd_sima_adduser
    do_mpd_sima_perms
    do_mpd_sima_create_db
    ;;

  abort-upgrade|abort-remove|abort-deconfigure)
    ;;

  *)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.


# Automatically added by dh_python3
if command -v py3compile >/dev/null 2>&1; then
	py3compile -p mpd-sima 
fi
if command -v pypy3compile >/dev/null 2>&1; then
	pypy3compile -p mpd-sima  || true
fi

# End automatically added section
# Automatically added by dh_installinit/13.21
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -z "${DPKG_ROOT:-}" ] && [ -x "/etc/init.d/mpd-sima" ]; then
		update-rc.d mpd-sima defaults >/dev/null
		if [ -n "$2" ]; then
			_dh_action=restart
		else
			_dh_action=start
		fi
		invoke-rc.d --skip-systemd-native mpd-sima $_dh_action || exit 1
	fi
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.21
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	# The following line should be removed in trixie or trixie+1
	deb-systemd-helper unmask 'mpd-sima.service' >/dev/null || true

	# was-enabled defaults to true, so new installations run enable.
	if deb-systemd-helper --quiet was-enabled 'mpd-sima.service'; then
		# Enables the unit on first installation, creates new
		# symlinks on upgrades if the unit file has changed.
		deb-systemd-helper enable 'mpd-sima.service' >/dev/null || true
	else
		# Update the statefile to add new symlinks (if any), which need to be
		# cleaned up on purge. Also remove old symlinks.
		deb-systemd-helper update-state 'mpd-sima.service' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.21
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -d /run/systemd/system ]; then
		systemctl --system daemon-reload >/dev/null || true
		if [ -n "$2" ]; then
			_dh_action=restart
		else
			_dh_action=start
		fi
		deb-systemd-invoke $_dh_action 'mpd-sima.service' >/dev/null || true
	fi
fi
# End automatically added section


exit 0
