# excalibur/patchfolder/roundcube-core/usr/lib/systemd/system/roundcube-gc.timer
[Unit]
Description=Purge expired Roundcube sessions, caches and tempfiles

[Timer]
OnCalendar=*-*-* *:05,35:00
Persistent=false

[Install]
WantedBy=timers.target
