# excalibur/patchfolder/bip/usr/lib/systemd/system/bip-config.service.d/bip-config.conf
[Service]
ExecStart=/usr/lib/systemd/scripts/bip_env.sh
Environment='BIP_DEFAULT_CONFIG=/etc/default/bip'
