# excalibur/patchfolder/xfconf/usr/lib/systemd/user/xfconfd.service
[Unit]
Description=Xfce configuration service

[Service]
Type=dbus
BusName=org.xfce.Xfconf
ExecStart=/usr/lib/x86_64-linux-gnu/xfce4/xfconf/xfconfd
