# excalibur/patchfolder/sitesummary/usr/lib/systemd/system/sitesummary.service
[Unit]
Description=Regular sitesummary maintenance

[Service]
Type=oneshot
ExecStart=/usr/lib/sitesummary/sitesummary-maintenance
