# excalibur/patchfolder/sitesummary/usr/lib/systemd/system/sitesummary.timer
[Unit]
Description=Timer for regular sitesummary maintenance

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
