# excalibur/patchfolder/debian-edu-config/usr/lib/systemd/system/debian-edu-update-netblock.service
[Unit]
Description=Update netblock according to netblock-hosts netgroup configuration.

[Service]
Type=oneshot
ExecStart=/usr/sbin/debian-edu-update-netblock auto
