# excalibur/patchfolder/dde-calendar/usr/lib/systemd/user/dde-calendar.service
[Unit]
Description = dde-calendar date service

[Service]
ExecStart=/usr/bin/dde-calendar
Type=dbus
BusName=com.deepin.Calendar
Slice=app.slice
