# excalibur/patchfolder/ipip/usr/lib/systemd/system/ipip.service
[Unit]
Description=IP over IP encapsulation daemon

[Service]
Type=simple
ExecStart=/usr/sbin/ipip

[Install]
WantedBy=multi-user.target

