# excalibur/patchfolder/leafnode/usr/lib/systemd/system/leafnode-fetchnews.service
[Unit]
Description=leafnode: Fetch news
After=network-online.target

[Service]
Type=oneshot
User=news
ExecStart=/usr/sbin/fetchnews
