# excalibur/patchfolder/leafnode/usr/lib/systemd/system/leafnode@.service
[Unit]
Description=NNTP server for small sites
Documentation=man:leafnode(8)

[Service]
ExecStart=/usr/sbin/leafnode
StandardInput=socket

