# excalibur/patchfolder/leafnode/usr/lib/systemd/system/leafnode.socket
[Unit]
Description=NNTP server for small sites (socket)

[Socket]
ListenStream=119
Accept=yes

[Install]
WantedBy=sockets.target
