# excalibur/patchfolder/postfix/usr/lib/systemd/system/postfix-resolvconf.path
[Unit]
Description=Watch for resolv.conf updates for postfix
ConditionPathExists=/etc/resolv.conf

[Path]
PathChanged=/etc/resolv.conf

[Install]
WantedBy=multi-user.target
