# excalibur/patchfolder/contactsd/usr/lib/systemd/user/contactsd.service
[Unit]
Description=Contacts daemon
After=lomiri.service
Requires=dbus.socket

[Service]
ExecStart=/usr/bin/contactsd
Restart=always

[Install]
WantedBy=lomiri.service
