# excalibur/patchfolder/gpaste-2/usr/lib/systemd/user/org.gnome.GPaste.Preferences.service
[Unit]
Description=GPaste preferences

[Service]
Type=dbus
BusName=org.gnome.GPaste.Preferences
ExecStart=/usr/libexec/gpaste/gpaste-preferences --gapplication-service
