# excalibur/patchfolder/oar-node/usr/lib/systemd/system/oar-nodecheck.service
[Unit]
Description=OAR nodecheck service
Documentation=man:oarnodecheckrun(8)
After=oar-node.service

[Service]
User=root
ExecStart=/usr/sbin/oarnodecheckrun
