# excalibur/patchfolder/atop/usr/lib/systemd/system/atop-rotate.service
[Unit]
Description=Restart atop daemon to rotate logs
Documentation=man:atop(1)

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart atop.service
