# excalibur/patchfolder/kylin-process-manager-daemon/usr/lib/systemd/system/kylin-process-manager-daemon.service
[Unit]
Description=kylin process manager daemon
After=display-manager.service

[Service]
Type=dbus
BusName=com.kylin.ProcessManager
ExecStart=/usr/bin/kylin-process-manager-daemon

[Install]
WantedBy=multi-user.target