# excalibur/patchfolder/intel-hdcp/usr/lib/systemd/system/hdcpd.service
[Unit]
Description=hdcp daemon

[Service]
Type=simple
ExecStart=/usr/bin/hdcpd

[Install]
WantedBy=multi-user.target
