# excalibur/patchfolder/aodh-evaluator/usr/lib/systemd/system/aodh-evaluator.service
[Unit]
Description=OpenStack Aodh Evaluator (aodh-evaluator)
After=mongodb.service rabbitmq-server.service frr.service network-online.target local-fs.target remote-fs.target 
Wants=mongodb.service rabbitmq-server.service frr.service network-online.target 

Documentation=man:aodh-evaluator(1)

[Service]
User=aodh
Group=aodh
Type=simple
WorkingDirectory=~
RuntimeDirectory=aodh lock/aodh
CacheDirectory=aodh
ExecStart=/etc/init.d/aodh-evaluator systemd-start
Restart=on-failure
LimitNOFILE=65535
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
