# excalibur/patchfolder/vitrage-api/usr/lib/systemd/system/vitrage-api.service
[Unit]
Description=OpenStack Vitrage RCAaaS API
After=postgresql.service mysql.service keystone.service rabbitmq-server.service ntp.service network-online.target local-fs.target remote-fs.target 
Wants=postgresql.service mysql.service keystone.service rabbitmq-server.service ntp.service network-online.target 

Documentation=man:vitrage-api(1)

[Service]
User=vitrage
Group=vitrage
Type=simple
WorkingDirectory=~
RuntimeDirectory=vitrage lock/vitrage
CacheDirectory=vitrage
ExecStart=/etc/init.d/vitrage-api systemd-start
Restart=on-failure
LimitNOFILE=65535
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
