# excalibur/patchfolder/deluged/DEBIAN/postinst
#!/bin/sh

set -e

DELUGE_USER="debian-deluged"
DELUGE_GROUP="debian-deluged"

Add_override ()
{
	if ! dpkg-statoverride --list ${4} > /dev/null 2>&1
	then
		dpkg-statoverride --update --add ${@}
	fi
}

case "${1}" in
	configure)
		if ! getent group ${DELUGE_GROUP} > /dev/null 2>&1
		then
			groupadd --system ${DELUGE_GROUP}
		fi

		if ! getent passwd ${DELUGE_USER} > /dev/null 2>&1
		then
			useradd --system -g ${DELUGE_GROUP} --home-dir /var/lib/deluged ${DELUGE_USER}
		fi

		Add_override ${DELUGE_USER} adm 2750 /var/log/deluged
		Add_override ${DELUGE_USER} ${DELUGE_GROUP} 0750 /var/lib/deluged
		Add_override ${DELUGE_USER} ${DELUGE_GROUP} 0750 /var/lib/deluged/config
		;;

	abort-upgrade|abort-remove|abort-deconfigure)

		;;

	*)
		echo "postinst called with unknown argument \`${1}'" >&2
		exit 1
		;;
esac

# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	# The following line should be removed in trixie or trixie+1
	deb-systemd-helper unmask 'deluged.service' >/dev/null || true

	# was-enabled defaults to true, so new installations run enable.
	if deb-systemd-helper --quiet was-enabled 'deluged.service'; then
		# Enables the unit on first installation, creates new
		# symlinks on upgrades if the unit file has changed.
		deb-systemd-helper enable 'deluged.service' >/dev/null || true
	else
		# Update the statefile to add new symlinks (if any), which need to be
		# cleaned up on purge. Also remove old symlinks.
		deb-systemd-helper update-state 'deluged.service' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -d /run/systemd/system ]; then
		systemctl --system daemon-reload >/dev/null || true
		if [ -n "$2" ]; then
			_dh_action=restart
		else
			_dh_action=start
		fi
		deb-systemd-invoke $_dh_action 'deluged.service' >/dev/null || true
	fi
fi
# End automatically added section


exit 0
