# excalibur/patchfolder/openqa/usr/lib/systemd/system/openqa-minion-restart.path
[Path]
PathChanged=/usr/share/perl5/Minion/Backend/resources/migrations/pg.sql

[Install]
WantedBy=multi-user.target
