# excalibur/patchfolder/pure-ftpd-ldap/DEBIAN/prerm
#! /bin/sh
# prerm script for pure-ftpd

set -e

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see /usr/share/doc/packaging-manual/

case "$1" in
	remove)
	update-inetd --pattern '/usr/sbin/pure-ftpd' --remove ftp
	;;
    upgrade|deconfigure)
	update-inetd --pattern '/usr/sbin/pure-ftpd' --disable ftp
    ;;
    failed-upgrade)
    ;;
    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 0
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

# Automatically added by dh_installinit/13.23
if [ -z "${DPKG_ROOT:-}" ] && [ "$1" = remove ] && [ -x "/etc/init.d/pure-ftpd-ldap" ] ; then
	invoke-rc.d pure-ftpd-ldap stop || exit 1
fi
# End automatically added section
# Automatically added by dh_installdeb/13.23
dpkg-maintscript-helper symlink_to_dir /usr/share/doc/pure-ftpd-ldap pure-ftpd-common 1.0.49-4.1\~ -- "$@"
# End automatically added section


exit 0
