# excalibur/patchfolder/geoipupdate/usr/lib/systemd/system/geoipupdate.timer
[Unit]
Description=Periodic GeoIP update
Documentation=man:geoipupdate(1)

[Timer]
OnCalendar=Wed,Sat America/New_York
RandomizedDelaySec=3h
Persistent=true

[Install]
WantedBy=timers.target
