# excalibur/patchfolder/lxcfs/etc/sv/lxcfs/finish
#!/bin/sh

fusermount -u /var/lib/lxcfs 2> /dev/null

exec /usr/lib/runit/finish-exec "$@"
