# excalibur/patchfolder/runit-services/usr/share/runit/sv/uuidd/run
#!/usr/bin/env /lib/runit/invoke-run
#Copyright: 2022 Cloux <cloux@rote.ch>
#            2023 Lorenzo Puliti <plorenzo@disroot.org>
#License: CC0-1.0
#Comment: adapted to Debian from https://github.com/cloux/runit-base

exec 2>&1

[ -d /run/uuidd ] || install -o uuidd -g uuidd -d /run/uuidd

exec /usr/sbin/uuidd --no-fork --no-pid
