# excalibur/patchfolder/runit-services/usr/share/runit/sv/zcfan/run
#!/usr/bin/env /lib/runit/invoke-run
#Copyright: 2023 Martin Steigerwald <martin@lichtvoll.de>
#License: BSD-3-Clause

exec 2>&1

exec /usr/bin/zcfan
