# excalibur/patchfolder/fake-hwclock/usr/lib/systemd/system/fake-hwclock-save.timer
[Unit]
Description=Periodically save current clock

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
