# excalibur/patchfolder/socklog-run/etc/sv/socklog-unix/run
#!/bin/sh
exec 2>&1
exec chpst -U _socklog-unix socklog unix /dev/log
