# excalibur/patchfolder/entropybroker/usr/lib/systemd/system/eb_client_file.service
[Unit]
Description=eb_client_file
Documentation=man:eb_client_file(1)

[Service]
Type=forking
EnvironmentFile=/etc/entropybroker/eb_client_file.config
ExecStart=/usr/bin/eb_client_file ${ARG01} ${ARG02} ${ARG03} ${ARG04} ${ARG05} ${ARG06} ${ARG07} ${ARG08} ${ARG09} ${ARG10} ${ARG11} ${ARG12} ${ARG13} ${ARG14} ${ARG15} ${ARG16}

[Install]
WantedBy=multi-user.target
