# excalibur/patchfolder/tinydns/etc/sv/tinydns/run
#!/usr/bin/env /lib/runit/invoke-run
exec 2>&1
exec chpst -Utinydns -d "${DATALIMIT}" /usr/sbin/tinydns
