# excalibur/patchfolder/privoxy/DEBIAN/postrm
#! /bin/sh
# postrm script for privoxy
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postrm> `remove'
#        * <postrm> `purge'
#        * <old-postrm> `upgrade' <new-version>
#        * <new-postrm> `failed-upgrade' <old-version>
#        * <new-postrm> `abort-install'
#        * <new-postrm> `abort-install' <old-version>
#        * <new-postrm> `abort-upgrade' <old-version>
#        * <disappearer's-postrm> `disappear' <r>overwrit>r> <new-version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    purge)
        # deluser --quiet privoxy
        for ext in .ucf-new .ucf-old .ucf-dist ""; do
                rm -f "/etc/privoxy/config$ext"
        done
        if [ -x "$(command -v ucf)" ]; then
            ucf --purge /etc/privoxy/config
        fi
        if [ -x "$(command -v ucfr)" ]; then
            ucfr --purge privoxy /etc/privoxy/config
        fi
        rm -rf /var/log/privoxy
        if [ -e /usr/share/debconf/confmodule ]; then
            . /usr/share/debconf/confmodule
            db_purge
        fi
	# Remove https inspection certificates on purge:
	if [ -d /var/lib/privoxy/certs ]; then
	    rm -rf /var/lib/privoxy/certs
	fi
        ;;

    remove|upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
        ;;

    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1

esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

# Automatically added by dh_installinit/13.24.1
if [ "$1" = "remove" ] && [ -x "/etc/init.d/privoxy" ] ; then
	chmod -x "/etc/init.d/privoxy" >/dev/null || true
fi
if [ -z "$DPKG_ROOT" ] && [ "$1" = "purge" ] ; then
	update-rc.d privoxy remove >/dev/null
fi
# End automatically added section
# Automatically added by dh_installdeb/13.24.1
dpkg-maintscript-helper rm_conffile /etc/privoxy/global.action 3.0.33-4\~ privoxy -- "$@"
dpkg-maintscript-helper rm_conffile /etc/privoxy/standard.action 3.0.33-4\~ privoxy -- "$@"
dpkg-maintscript-helper rm_conffile /etc/privoxy/templates/show-version 3.0.33-4\~ privoxy -- "$@"
dpkg-maintscript-helper rm_conffile /etc/privoxy/templates/edit-actions-for-url-string-filter 3.0.33-3\~ privoxy -- "$@"
# End automatically added section
# Automatically added by dh_installsystemd/13.24.1
if [ "$1" = remove ] && [ -d /run/systemd/system ] ; then
	systemctl --system daemon-reload >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.1
if [ "$1" = "purge" ]; then
	if [ -x "/usr/bin/deb-systemd-helper" ]; then
		deb-systemd-helper purge 'privoxy-cleanup.timer' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.1
if [ "$1" = remove ] && [ -d /run/systemd/system ] ; then
	systemctl --system daemon-reload >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.1
if [ "$1" = "purge" ]; then
	if [ -x "/usr/bin/deb-systemd-helper" ]; then
		deb-systemd-helper purge 'privoxy.service' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_installdebconf/13.24.1
if [ "$1" = purge ] && [ -e /usr/share/debconf/confmodule ]; then
	. /usr/share/debconf/confmodule
	db_purge
fi
# End automatically added section
# Automatically added by dh_apparmor/3.1.7-1
if [ "$1" = "purge" ] && ! [ -e "/etc/apparmor.d/usr.sbin.privoxy" ] ; then
    rm -f "/etc/apparmor.d/disable/usr.sbin.privoxy" || true
    rm -f "/etc/apparmor.d/force-complain/usr.sbin.privoxy" || true
    rm -f "/etc/apparmor.d/local/usr.sbin.privoxy" || true
    rm -f /var/cache/apparmor/*/"usr.sbin.privoxy" || true
    rmdir /etc/apparmor.d/disable 2>/dev/null || true
    rmdir /etc/apparmor.d/local   2>/dev/null || true
    rmdir /etc/apparmor.d         2>/dev/null || true
fi
# End automatically added section


exit 0
