# excalibur/patchfolder/runit-services/usr/share/runit/sv/gdomap/run
#!/usr/bin/env /lib/runit/invoke-run
#Copyright: 2022 Lorenzo Puliti <plorenzo@disroot.org>
#License: CC0-1.0

sv start default-syslog ||  true

if [ ! -d /run/gdomap ]; then
        install -d -m 0755 /run/gdomap
fi

exec 2>&1

exec /usr/bin/gdomap  -f $GDOMAP_ARGS
