# excalibur/patchfolder/guerillabackup/DEBIAN/postinst
#!/bin/sh
# postinst script for guerillabackup
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.


# Automatically added by dh_python3
if command -v py3compile >/dev/null 2>&1; then
	py3compile -p guerillabackup /usr/lib/guerillabackup
fi
if command -v pypy3compile >/dev/null 2>&1; then
	pypy3compile -p guerillabackup /usr/lib/guerillabackup || true
fi

# End automatically added section
# Automatically added by dh_installsystemd/13.20
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'guerillabackup-generator.service'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'guerillabackup-generator.service' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'guerillabackup-generator.service'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'guerillabackup-generator.service' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'guerillabackup-generator.service' >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.20
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed 'guerillabackup-transfer.service'; then
		# The following line should be removed in trixie or trixie+1
		deb-systemd-helper unmask 'guerillabackup-transfer.service' >/dev/null || true

		if deb-systemd-helper --quiet was-enabled 'guerillabackup-transfer.service'; then
			# Create new symlinks, if any.
			deb-systemd-helper enable 'guerillabackup-transfer.service' >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state 'guerillabackup-transfer.service' >/dev/null || true
fi
# End automatically added section


# Start the service only when it was already enabled before updating.
# This is required as "dh_systemd_start" is disabled in rules
# file, thus not restarting the services. See rules file for more
# information.
if test -d /run/systemd/system && test -e /run/guerillabackup.dpkg-update.run-state; then
  daemonReloadedFlag="false"
  for serviceName in guerillabackup-generator.service guerillabackup-transfer.service; do
    if grep -q -e "${serviceName}: active" -- /run/guerillabackup.dpkg-update.run-state; then
      if [ "${daemonReloadedFlag}" != "true" ]; then
        systemctl --system daemon-reload >/dev/null || true
        daemonReloadedFlag="true"
      fi
      deb-systemd-invoke restart "${serviceName}" >/dev/null || true
    fi
  done
fi
rm -f -- /run/guerillabackup.dpkg-update.run-state
exit 0
