# excalibur/patchfolder/usbguard-notifier/usr/lib/systemd/user/usbguard-notifier.service
[Unit]
Description=USBGuard Notifier
After=usbguard.service

[Service]
ExecStart=/usr/bin/usbguard-notifier

[Install]
WantedBy=default.target
