# excalibur/patchfolder/runit-services/usr/share/runit/sv/lircd/run
#!/usr/bin/env /lib/runit/invoke-run
#Copyright: 2022 Lorenzo Puliti <plorenzo@disroot.org>
#License: CC0-1.0

if [ ! -f /etc/lirc/lirc_options.conf ]; then
	exit 162
fi
exec 2>&1

/usr/sbin/lircd-setup
[ -d "/run/lirc" ] || mkdir -p "/run/lirc"

exec /usr/sbin/lircd --nodaemon --logfile /dev/stdout
