# excalibur/patchfolder/qcom-phone-utils/usr/lib/systemd/system/qcom-modem-setup.service
[Unit]
Description=Qualcomm modem setup
Before=ModemManager.service

[Service]
Type=oneshot
ExecStart=/usr/bin/qcom-modem-setup

[Install]
WantedBy=ModemManager.service
