# excalibur/patchfolder/designate-pool-manager/usr/lib/systemd/system/designate-pool-manager.service
[Unit]
Description=OpenStack Designate DNSaaS pool manager (designate-pool-manager)
After=postgresql.service mysql.service keystone.service rabbitmq-server.service ntp.service network-online.target local-fs.target remote-fs.target 
Wants=postgresql.service mysql.service keystone.service rabbitmq-server.service ntp.service network-online.target 

Documentation=man:designate-pool-manager(1)

[Service]
User=designate
Group=designate
Type=simple
WorkingDirectory=~
RuntimeDirectory=designate lock/designate
CacheDirectory=designate
ExecStart=/etc/init.d/designate-pool-manager systemd-start
Restart=on-failure
LimitNOFILE=65535
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
