# excalibur/patchfolder/runit-services/usr/share/runit/sv/sddm/run
#!/usr/bin/env /lib/runit/invoke-run
#Copyright: 2022 Lorenzo Puliti <plorenzo@disroot.org>
#License: CC0-1.0

DEFAULT_DISPLAY_MANAGER=$(cat /etc/X11/default-display-manager)

# don't start if sddm is not the default DM
if [ $DEFAULT_DISPLAY_MANAGER != /usr/bin/sddm ]; then
	exit 161
fi

sv start dbus || true
sv start elogind || true

mkdir -p /var/lib/sddm
chown -R sddm:sddm /var/lib/sddm

exec 2>&1

exec /usr/bin/sddm
