# excalibur/patchfolder/proxsmtp/usr/lib/systemd/system/proxsmtp.target
[Unit]
Description=multi purpose SMTP proxy

[Install]
WantedBy=multi-user.target
